# Deploy/conf/certbot_renew.service
[Unit]
Description=Certbot renew

[Service]
ExecStart=/bin/bash /root/LegacyPlayersV3/Deploy/renew_cert.sh
User=root