# Deploy/conf/restart_backend.service
[Unit]
Description=Restart Backend Service

[Service]
ExecStart=/bin/bash /root/LegacyPlayersV3/Deploy/restart_backend.sh
User=root